# com.whisper-mac.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.whisper-mac</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>APP_PATH/start.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>APP_PATH/whisper-mac.log</string>
    <key>StandardErrorPath</key>
    <string>APP_PATH/whisper-mac-error.log</string>
</dict>
</plist> 